# packages/opensnitch/daemon-deb/.post/install.sh
# [[[cog
# ]]]
# [[[end]]]

printf "You might need to enable the service if it has not been enabled automatically\n";
printf "\tsudo systemctl enable opensnitch; sudo systemctl start opensnitch\n";